# systemd/edu.service
[Unit]
After=network.target
Description=Docker project for edu.xiver.ru

[Service]
Type=simple
ExecStart=/opt/xiver/testing-system/scripts/start.sh
ExecStop=/opt/xiver/testing-system/scripts/stop.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target